# scripts/init.d/rainbow-monitoragent-alpine
#!/sbin/openrc-run
name="Rainbow Monitor Agent"         
description="Exposes via a webservice the core vitals of the machine"           
command="/usr/local/bin/rainbow-monitoragent-alpine.sh"
command_background=true
pidfile="/run/rainbow-monitoragent.pid"
